# systemd/system/arp.service
[Unit]
Description=arp
After=network.target

[Service]
ExecStart=/etc/scripts/staticarp.sh

[Install]
WantedBy=multi-user.target
