# cron.d/nbox
* * * * * root python /var/ntop/utils/scheduler.py > /dev/null 2>&1
